# apps/voice-satellite/satellite.service
[Unit]
Description=Mauricio Voice Satellite
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/mauricio-satellite
EnvironmentFile=/home/%i/mauricio-satellite/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/%i/mauricio-satellite/venv/bin/python satellite.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
